# provisioning/roles/demo/files/konk-nymph.service
[Unit]
Description=Konk node local daemon service.

[Service]
Type=simple
ExecStart=/usr/bin/konk --config /etc/konk.yaml nymph -vvvvv
KillMode=control-group

[Install]
WantedBy=multi-user.target